# srcs/mySql/srcs/supervisord.conf
[supervisord]
nodaemon=true

[program:mysql]
command=sh run.sh
autostart=true
autorestart=true