# srcs/grafana/srcs/supervisord.conf
[supervisord]
nodaemon=true

[program:grafana]
command=grafana-7.1.5/bin/ "./grafana-server"
autostart=true
autorestart=true


[program:telegraf]
command=telegraf
autostart=true
autorestart=true