# srcs/nginx/srcs/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=nginx -g 'daemon off;'
autostart=true
autorestart=true

[program:sshd]
command=sh run.sh
autorestart=true
autostart=true
nodaemon=true
daemon=off


[program:telegraf]
command=telegraf
autostart=true
autorestart=true