# srcs/nginx/srcs/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=nginx -g 'daemon off;'
autostart=true
autorestart=true

[program:telegraf]
command=telegraf
autostart=true
autorestart=true