# srcs/nginx/srcs/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=nginx -g 'daemon off;'
autostart=true
autorestart=true

#[program:ssh]
#command=/usr/sbin/sshd
#autorestart=true
#autostart=true